# packages/mediacenter/kodi/system.d/kodi.service
[Unit]
Description=Kodi Media Center
After=graphical.target
Requires=graphical.target

[Service]
Environment=HOME=/storage
ExecStart=/usr/lib/kodi/kodi.sh
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=5
Restart=always
RestartSec=2
StartLimitInterval=0
LimitNOFILE=16384

[Install]
WantedBy=kodi.target
