# packages/base/busybox/system.d/autostart.service
[Unit]
Description=Autostart script
ConditionPathExists=/storage/.config/autostart.sh

[Service]
Type=oneshot
ExecStart=-/bin/sh -c ". /etc/profile; exec /bin/sh /storage/.config/autostart.sh"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
