# packages/sysutils/busybox/system.d/hostname.service
[Unit]
Description=Setup Hostname
Before=network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/openelec/hostname-setup

[Install]
WantedBy=multi-user.target
